# deployment/services.prod.ini
[supervisord]
loglevel=debug
nodaemon = true
user = root

[program:gunicorn]
command = gunicorn basechem.wsgi:application --bind 0.0.0.0:8000 --timeout 300 --log-level=debug
autostart = True
autorestart = True
directory = /home/app/web
user = root

[program:django-q]
command = python /home/app/web/manage.py qcluster 
directory = /home/app/web
stopasgroup = true
user = root

[program:django-q-fast]
environment=Q_CLUSTER_NAME=fast
command = python /home/app/web/manage.py qcluster
directory = /home/app/web
stopasgroup = true
user = root

[program:django-q-slow]
environment=Q_CLUSTER_NAME=slow
command = python /home/app/web/manage.py qcluster
directory = /home/app/web
stopasgroup = true
user = root
